# etc/supervisor.d/domain2.ini
[program:domain2]
user = nginx
directory = /opt/dummy/dummy
command = /opt/py_368/bin/gunicorn --workers 3 -b unix:/run/gunicorn/domain2.socket
autostart = yes
autorestart = yes
startretries = 3
stdout_logfile = /opt/dummy/log/%(program_name)s.log
stderr_logfile = /opt/dummy/log/%(program_name)s.log
environment=FLASK_APP="app:create_app('domain2')"
